# systemd/analyzers/conversation-analyzer.service
[Unit]
Description=JIB Conversation Analyzer
Documentation=https://github.com/jwiesebron/james-in-a-box
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%u

# Execute via jib (runs analyzer in container environment)
# Note: Requires jib script to have --exec flag implemented
ExecStart=%h/khan/james-in-a-box/jib --exec /home/jwies/tools/analyze-conversations --days 7

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conversation-analyzer

# Restart behavior
Restart=no

[Install]
WantedBy=multi-user.target
