# host-services/analysis/github-watcher/comment-responder.timer
[Unit]
Description=Comment Responder Timer - Monitor every 5 minutes
Documentation=file://%h/khan/james-in-a-box/host-services/analysis/github-watcher/README.md

[Timer]
# Run every 5 minutes based on when the timer was last activated
# No OnBootSec to avoid immediate execution during setup
OnUnitActiveSec=5min
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
