# host-services/analysis/github-watcher/github-watcher.service
[Unit]
Description=GitHub Watcher - Monitor GitHub and trigger jib analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/khan/james-in-a-box/host-services/analysis/github-watcher/github-watcher.py
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin:%h/khan/james-in-a-box/bin

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
