# host-services/analysis/github-watcher/github-watcher.timer
[Unit]
Description=GitHub Watcher Timer - Monitor every 5 minutes
Documentation=file://%h/khan/james-in-a-box/host-services/analysis/github-watcher/README.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
