# components/service-monitor/service-failure-notify@.service
[Unit]
Description=Send Slack notification for failed service %i
Documentation=file:///home/jwies/khan/james-in-a-box/docs/reference/service-failure-notifications.md

[Service]
Type=oneshot
ExecStart=/home/jwies/khan/james-in-a-box/components/service-monitor/notify-service-failure.sh %i
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-failure-notify

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/jwies/sharing/notifications

[Install]
WantedBy=multi-user.target
