# host-services/utilities/worktree-watcher/worktree-watcher.service
[Unit]
Description=JIB Worktree Watcher - Clean up orphaned git worktrees
Documentation=file://%h/repos/james-in-a-box/host-services/utilities/worktree-watcher/README.md
# No network needed - local git worktree operations only
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot
ExecStart=%h/repos/james-in-a-box/host-services/utilities/worktree-watcher/worktree-watcher.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=worktree-watcher

# Security (baseline)
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
