# host-services/utilities/worktree-watcher/worktree-watcher.service
[Unit]
Description=JIB Worktree Watcher - Clean up orphaned git worktrees
Documentation=file://%h/khan/james-in-a-box/host-services/utilities/worktree-watcher/README.md
After=network.target

[Service]
Type=oneshot
ExecStart=%h/khan/james-in-a-box/host-services/utilities/worktree-watcher/worktree-watcher.sh
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
