# host-services/utilities/worktree-watcher/worktree-watcher.timer
[Unit]
Description=JIB Worktree Watcher Timer - Periodic cleanup of orphaned worktrees
Documentation=file://%h/khan/james-in-a-box/host-services/utilities/worktree-watcher/README.md

[Timer]
# Run 5 minutes after boot
OnBootSec=5min

# Run every 15 minutes
OnUnitActiveSec=15min

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
