# host-services/utilities/worktree-watcher/worktree-watcher.timer
[Unit]
Description=JIB Worktree Watcher Timer - Periodic cleanup of orphaned worktrees
Documentation=file://%h/khan/james-in-a-box/host-services/utilities/worktree-watcher/README.md

[Timer]
# Run every 15 minutes based on when the timer was last activated
# No OnBootSec to avoid immediate execution during setup
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
