# host-services/codebase-analyzer/codebase-analyzer.service
[Unit]
Description=Codebase Improvement Analyzer
Documentation=https://github.com/jwbron/james-in-a-box
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot

# Working directory
WorkingDirectory=/home/jwies/khan/james-in-a-box/host-services/codebase-analyzer

# Execute the analyzer via jib --exec (runs in ephemeral container with Claude Code)
# --implement: automatically fix issues and open a PR
# --max-fixes 10: limit to top 10 highest-priority fixes
ExecStart=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/host-services/codebase-analyzer/codebase-analyzer.py --implement --max-fixes 10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codebase-analyzer

# Restart behavior
Restart=no

[Install]
WantedBy=multi-user.target
