# host-services/sync/context-sync/context-sync.service
[Unit]
Description=Context Sync - Sync all connectors (Confluence, JIRA, etc.)
Documentation=file://%h/khan/james-in-a-box/host-services/sync/context-sync/README.md
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/khan/james-in-a-box/host-services/sync/context-sync
Environment="PATH=%h/khan/james-in-a-box/host-services/.venv/bin:/usr/local/bin:/usr/bin:/bin"
# Configuration is loaded from ~/.config/context-sync/.env by context-sync.py
# Uses shared virtual environment installed by setup.sh via uv
ExecStart=%h/khan/james-in-a-box/host-services/.venv/bin/python %h/khan/james-in-a-box/host-services/sync/context-sync/context-sync.py

# After syncing context data, run analysis in JIB container (new ephemeral containers)
# Confluence and JIRA processors analyze docs and tickets for updates
ExecStartPost=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/jib-container/jib-tasks/confluence/confluence-processor.py
ExecStartPost=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/jib-container/jib-tasks/jira/jira-processor.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=context-sync

# Resource limits
Nice=10
IOSchedulingClass=2
IOSchedulingPriority=7

# Security (baseline)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target

