# host-services/analysis/adr-researcher/adr-researcher.service
[Unit]
Description=ADR Researcher - Research and update ADRs with industry findings
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/khan/james-in-a-box/host-services/analysis/adr-researcher/adr-researcher.py --scope merged
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin:%h/khan/james-in-a-box/bin

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
