# host-services/analysis/adr-researcher/adr-researcher.timer
[Unit]
Description=Weekly ADR Researcher Timer
Documentation=file://%h/khan/james-in-a-box/host-services/analysis/adr-researcher/README.md

[Timer]
OnCalendar=Mon *-*-* 11:00:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
