# host-services/analysis/adr-researcher/adr-researcher.timer
[Unit]
Description=Run ADR Researcher daily to update implemented ADRs

[Timer]
OnBootSec=15min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
