# host-services/analysis/feature-analyzer/feature-analyzer-weekly.timer
[Unit]
Description=Weekly feature analyzer run for FEATURES.md updates
Documentation=https://github.com/jwbron/james-in-a-box/blob/main/docs/adr/implemented/ADR-Feature-Analyzer-Documentation-Sync.md
Requires=feature-analyzer-weekly.service

[Timer]
# Runs Mondays at 11:00 AM (local time)
# Note: Uses system timezone (local time). For consistent behavior:
# - Ensure host is configured to America/Los_Angeles (Pacific Time) per ADR recommendations
# - Alternatively, use explicit UTC: OnCalendar=Mon *-*-* 18:00:00 UTC (PST) or 19:00:00 UTC (PDT)
# - Check current timezone with: timedatectl show --property=Timezone
OnCalendar=Mon *-*-* 11:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
