# host-services/slack/slack-notifier/slack-notifier.service
[Unit]
Description=Slack Notifier - Monitor Claude sandbox and send Slack notifications
Documentation=file://%h/khan/james-in-a-box/host-services/slack/slack-notifier/README.md
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=simple
WorkingDirectory=%h/khan/james-in-a-box
# Configuration loaded from ~/.config/jib-notifier/config.json
# Uses shared virtual environment installed by setup.sh via uv
ExecStart=%h/khan/james-in-a-box/host-services/.venv/bin/python %h/khan/james-in-a-box/host-services/slack/slack-notifier/slack-notifier.py

# Restart on failure
Restart=on-failure
RestartSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-notifier

# Resource limits
Nice=10
IOSchedulingClass=2
IOSchedulingPriority=7

# Security (baseline)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
