# host-services/conversation-analyzer/conversation-analyzer.service
[Unit]
Description=JIB Conversation Analyzer
Documentation=https://github.com/jwbron/james-in-a-box
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot

# Execute the conversation analyzer directly from the repo
ExecStart=/usr/bin/python3 %h/khan/james-in-a-box/components/conversation-analyzer/conversation-analyzer.py --days 7

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conversation-analyzer

# Restart behavior
Restart=no

[Install]
WantedBy=multi-user.target
