# components/github-sync/systemd/github-sync.timer
[Unit]
Description=GitHub PR Sync Timer - Runs every 15 minutes
Documentation=file://%h/khan/james-in-a-box/components/github-sync/README.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
