# host-services/utilities/github-token-refresher/github-token-refresher.timer
[Unit]
Description=GitHub Token Refresher Timer - Refresh tokens every 5 minutes
Documentation=file://%h/khan/james-in-a-box/host-services/utilities/github-token-refresher/README.md

[Timer]
# Run every 5 minutes
# GitHub App tokens expire after 1 hour, but shorter interval provides
# better reliability when system resumes from suspend
OnCalendar=*:0/5

# Run immediately on boot to ensure token is fresh
OnBootSec=10s

# CRITICAL: Run missed executions after system resume from suspend/hibernate
# This ensures tokens are refreshed immediately after waking from sleep
Persistent=true

# Small random delay to avoid thundering herd on multi-machine setups
RandomizedDelaySec=10s

[Install]
WantedBy=timers.target
