# host-services/analysis/beads-analyzer/beads-analyzer.service
[Unit]
Description=JIB Beads Integration Analyzer
Documentation=https://github.com/jwbron/james-in-a-box
After=network-online.target
Wants=network-online.target
# Use ConditionPathExists instead of Requires=docker.service
# User units cannot depend on system services
ConditionPathExists=/var/run/docker.sock
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot

# Execute the beads analyzer - delegates to jib container for isolated worktree
# The container-side processor handles all file operations to avoid tainting main worktree
ExecStart=/usr/bin/python3 %h/repos/james-in-a-box/host-services/analysis/beads-analyzer/beads-analyzer.py --days 7

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beads-analyzer

# Restart behavior
Restart=no

# Environment for jib
Environment="PATH=/usr/local/bin:/usr/bin:/bin:%h/repos/james-in-a-box/bin"

[Install]
WantedBy=multi-user.target
