# host-services/analysis/conversation-analyzer/conversation-analyzer.service
[Unit]
Description=JIB Conversation Analyzer - Analyze conversation patterns
Documentation=file://%h/khan/james-in-a-box/host-services/analysis/conversation-analyzer/README.md
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot

# Execute the conversation analyzer directly from the repo
ExecStart=/usr/bin/python3 %h/khan/james-in-a-box/host-services/analysis/conversation-analyzer/conversation-analyzer.py --days 7

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conversation-analyzer

# Restart behavior (oneshot - no restart)
Restart=no

# Security (baseline)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
