# components/codebase-analyzer/codebase-analyzer.service
[Unit]
Description=Codebase Improvement Analyzer
Documentation=https://github.com/jwiesebron/james-in-a-box
After=network-online.target
Wants=network-online.target
OnFailure=service-failure-notify@%n.service

[Service]
Type=oneshot

# Working directory
WorkingDirectory=/home/jwies/khan/james-in-a-box/components/codebase-analyzer

# Execute the analyzer via jib --exec (runs in ephemeral container with Claude Code)
ExecStart=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/components/codebase-analyzer/codebase-analyzer.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codebase-analyzer

# Restart behavior
Restart=no

[Install]
WantedBy=multi-user.target
