# scripts/codebase-analyzer.service
[Unit]
Description=Codebase Improvement Analyzer
Documentation=https://github.com/jwiesebron/james-in-a-box
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jwies
Group=jwies

# Working directory
WorkingDirectory=/home/jwies/khan/james-in-a-box/scripts

# Execute the analyzer
ExecStart=/usr/bin/python3 /home/jwies/khan/james-in-a-box/scripts/codebase-analyzer.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codebase-analyzer

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jwies/sharing/notifications

# Restart behavior
Restart=no

[Install]
WantedBy=multi-user.target
