# host-services/venv-setup.service
[Unit]
Description=JIB Host Services - Python venv setup
Documentation=file://%h/repos/james-in-a-box/host-services/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/repos/james-in-a-box
ExecStart=/usr/bin/make venv

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=venv-setup

[Install]
WantedBy=default.target
