# host-services/sync/github-sync/github-sync.timer
[Unit]
Description=GitHub PR Sync Timer - Runs every 5 minutes
Documentation=file://%h/khan/james-in-a-box/host-services/github-sync/README.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
