# host-services/sync/github-sync/github-sync.service
[Unit]
Description=GitHub PR Sync - Sync open PRs and check status
Documentation=file://%h/khan/james-in-a-box/host-services/sync/github-sync/README.md

[Service]
Type=oneshot
WorkingDirectory=%h/khan/james-in-a-box/host-services/sync/github-sync

# Sync repos defined in config/repositories.yaml (the single source of truth)
# Use --use-config to load repo list and sync settings from the config file
ExecStart=%h/khan/james-in-a-box/host-services/sync/github-sync/github-sync.py --use-config

# After syncing PR data, run analysis in JIB container (new ephemeral container)
# 1. Issue fixer - detects check failures and merge conflicts, lets Claude determine fixes
ExecStartPost=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/jib-container/jib-tasks/github/issue-fixer.py

# 2. PR reviewer - auto-reviews new PRs from others (skips your own)
ExecStartPost=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/jib-container/jib-tasks/github/pr-reviewer.py --watch

# 3. Comment responder - detects new comments needing responses on YOUR PRs
ExecStartPost=%h/khan/james-in-a-box/bin/jib --exec python3 %h/khan/james-in-a-box/jib-container/jib-tasks/github/comment-responder.py

[Install]
WantedBy=default.target
